# build_files/kyth-default-flatpaks.service
[Unit]
Description=Install game-ready Flatpak applications
Documentation=https://github.com/mrtrick37/kyth
After=network-online.target kyth-flathub-setup.service
# Requires= (not just Wants=) so this service is skipped on offline boots
# rather than starting, failing, and leaving the sentinel file unwritten.
Requires=network-online.target kyth-flathub-setup.service
# Sentinel is VERSIONED — bump the filename (v2 → v3, etc.) any time the app
# list changes. The old sentinel survives OS upgrades (it lives in /var which
# bootc never touches), so without versioning existing users never see new apps.
# Current version: v3 (added Bottles, Flatseal, WebApp Manager, Warpinator)
ConditionPathExists=!/var/lib/kyth/default-flatpaks-v3-done

[Service]
Type=oneshot
# --system  : install for all users (consistent with kyth-flathub-setup.service)
# --or-update: idempotent — updates if already present, installs if not
# Keep first boot game-ready without pulling in general desktop, creator, or
# admin utilities. Those stay available from the KythOS System Hub.
ExecStart=/usr/bin/flatpak install -y --noninteractive --system --or-update flathub \
    com.valvesoftware.Steam \
    net.lutris.Lutris \
    com.heroicgameslauncher.hgl \
    com.github.Matoking.protontricks \
    net.davidotek.pupgui2 \
    org.x.Warpinator \
    com.usebottles.bottles \
    com.github.tchx84.Flatseal \
    io.github.vikdevelop.WebApp
ExecStartPost=/bin/mkdir -p /var/lib/kyth
ExecStartPost=/bin/touch /var/lib/kyth/default-flatpaks-v3-done
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
