# build_files/kyth-flathub-setup.service
[Unit]
Description=Add Flathub remote for Flatpak
Documentation=https://github.com/mrtrick37/kyth
After=network-online.target
# Requires= (not just Wants=) so systemd does not start this service when
# there is no network — Wants= is advisory and allows the unit to start even
# if network-online never becomes active (e.g. offline boot).
Requires=network-online.target
# Flatpak remotes live in /var — only configurable at runtime, not build time.
# This service is idempotent: --if-not-exists means it exits immediately if
# Flathub is already configured (e.g. after the first successful run).
# Retry limit for DNS-not-ready races at first boot.
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
ExecStart=flatpak remote-add --if-not-exists --system flathub https://dl.flathub.org/repo/flathub.flatpakrepo
RemainAfterExit=yes
# Retry on DNS-not-ready races at first boot (network-online fires before
# resolver is actually usable on some systems/boot paths).
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
