# opt/paloaltonetworks/globalprotect/gpd.service
[Unit]
Description=GlobalProtect VPN client daemon
After=network.target

[Service]
Type=simple
ExecStartPre=/opt/paloaltonetworks/globalprotect/pre_exec_gps.sh
ExecStart=/opt/paloaltonetworks/globalprotect/PanGPS
Restart=on-failure
RestartSec=5
# Stop restarting after 3 failures in 30 s. Without a limit, a missing or
# corrupt PanGPS binary fills the journal and adds ~5 s of delay to every
# subsequent boot as systemd works through the restart backoff queue.
StartLimitBurst=3
StartLimitIntervalSec=30
WorkingDirectory=/opt/paloaltonetworks/globalprotect

[Install]
WantedBy=multi-user.target

